# loto-ia.service
[Unit]
Description=LOTO IA - Angular SSR Application
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/tienda-web-lotoAI-1
ExecStart=/usr/bin/npm run start:prod
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=4000

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=loto-ia

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/tienda-web-lotoAI-1

[Install]
WantedBy=multi-user.target 